# root/etc/services.d/pgstudio/run
#!/usr/bin/with-contenv sh
exec /usr/share/pgstudio/bin/catalina.sh run
