# root/etc/services.d/pgstudio/run
#!/usr/bin/with-contenv sh
exec /usr/local/pgstudio/bin/catalina.sh run
