# root/etc/services.d/jupyterhub/run
#!/usr/bin/with-contenv sh
exec /usr/local/bin/jupyterhub -f /etc/jupyterhub/jupyterhub_config.py
