# root/etc/services.d/rstudio/finish
#!/usr/bin/with-contenv sh
rstudio-server stop
